# install/HTFanControl.service
[Unit]
Description=HTFanControl
After=multi-user.target

[Service]
ExecStart=/opt/HTFanControl/HTFanControl

[Install]
WantedBy=multi-user.target